# root/etc/services.d/cron/run
#!/bin/execlineb -P

# Redirect stderr to stdout.
fdmove -c 2 1

# Wait until openbox is running
#if { s6-svwait -t 10000 -U /var/run/s6/services/openbox/ }
# Wait until guacamole is running
#if { s6-svwait -t 10000 -U /var/run/s6/services/guacamole/ }

if -n { test -f /var/run/crond.pid }
/usr/sbin/cron
